# deploy/kogu.service
[Unit]
Description=Kogu dictionary API (Rust/Axum, read-only SQLite)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/miro/kogu/backend
ExecStart=/srv/miro/kogu/bin/kogu
Restart=on-failure
RestartSec=2
Environment=KOGU_DB=/srv/miro/kogu/data/kogu.sqlite
Environment=KOGU_PORT=4100
Environment=RUST_LOG=kogu=info
# OCR: dlopen the official onnxruntime .so (load-dynamic) + PP-OCRv5 models, both on the volume
Environment=KOGU_OAR_DIR=/mnt/HC_Volume_102319212/wenbun/oar
Environment=ORT_DYLIB_PATH=/mnt/HC_Volume_102319212/wenbun/oar/libonnxruntime.so
# read-only service: reads the DB + OCR models, proxies handwriting ink
ProtectSystem=strict
ReadOnlyPaths=/srv/miro/kogu /mnt/HC_Volume_102319212/wenbun/oar
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
