# deploy/kogu.service
[Unit]
Description=Kogu dictionary API (Rust/Axum, read-only SQLite)
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/miro/kanzi/backend
ExecStart=/srv/miro/kanzi/bin/kogu
Restart=on-failure
RestartSec=2
Environment=KANZI_DB=/srv/miro/kanzi/data/kanzi.sqlite
Environment=KANZI_PORT=4100
Environment=RUST_LOG=kanzi=info
# read-only service: it only reads the DB + proxies handwriting ink
ProtectSystem=strict
ReadOnlyPaths=/srv/miro/kanzi
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
